# x11/freebsd-live-common-settings/pkg-plist
etc/polkit-1/localauthority/10-vendor.d/10-live-cd.pkla
share/freebsd/common-live-settings/base/backup/etc/gettytab
share/freebsd/common-live-settings/base/backup/etc/rc.d/mountcritlocal
share/freebsd/common-live-settings/base/backup/etc/ttys
share/freebsd/common-live-settings/base/override/etc/login.conf
share/freebsd/common-live-settings/base/override/etc/rc.d/mountcritlocal
share/freebsd/common-live-settings/base/override/root/cardDetect/XF86Config.compat
share/freebsd/common-live-settings/base/override/root/cardDetect/XF86Config.default
share/freebsd/common-live-settings/base/override/root/cardDetect/XF86Config.intel
share/freebsd/common-live-settings/base/override/root/cardDetect/XF86Config.modesetting
share/freebsd/common-live-settings/base/override/root/cardDetect/XF86Config.qemu
share/freebsd/common-live-settings/base/override/root/cardDetect/XF86Config.scfb
share/freebsd/common-live-settings/base/override/root/cardDetect/XF86Config.virtualbox
share/freebsd/common-live-settings/base/override/root/cardDetect/XF86Config.vmware
share/freebsd/common-live-settings/base/override/root/cardDetect/setuptouch.sh
share/freebsd/common-live-settings/base/override/root/cardDetect/x_config.sh
share/freebsd/common-live-settings/base/override/root/cardDetect/xorg.conf.virtualbox
share/freebsd/common-live-settings/base/override/root/functions.sh
share/freebsd/common-live-settings/base/override/root/iso_to_hd.sh
share/freebsd/common-live-settings/base/override/root/sysconfig.sh
share/freebsd/common-live-settings/base/override/root/sysutil.sh
share/freebsd/common-live-settings/config-live-settings
share/freebsd/common-live-settings/etc/rc.d/cheatcodes
share/freebsd/common-live-settings/etc/rc.d/unionfs
share/freebsd/common-live-settings/etc/rc.d/uzip
